# Examples/TypeScript/Containers/NodeJS-MongoDB-REST/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D
autostart=true
autorestart=true

[program:mongodb_rest]
command=node ./dist/src/server.js
directory=/app
autostart=true
autorestart=true
priority=10

[program:sshd]
command=/usr/sbin/sshd -D
autostart=true
autorestart=true
priority=30
stderr_logfile=/var/log/supervisor/sshd-stderr.log
stdout_logfile=/var/log/supervisor/sshd-stdout.log