# Examples/PHP/Containers/Laravel-MongoDB-REST/supervisord.conf
[supervisord]
nodaemon=true

[program:php-fpm]
command=php-fpm --nodaemonize
autostart=true
autorestart=true
priority=10

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
priority=20

[program:sshd]
command=/usr/sbin/sshd -D
autostart=true
autorestart=true
priority=30
stderr_logfile=/var/log/supervisor/sshd-stderr.log
stdout_logfile=/var/log/supervisor/sshd-stdout.log


;[program:queue-worker]
;command=php /var/www/html/artisan queue:work --tries=3
;autostart=true
;autorestart=true
